# roles/opencode/templates/opencode.service.j2
[Unit]
Description=OpenCode Web
After=network.target

[Service]
Type=simple
User={{ opencode_user }}
Group={{ opencode_group }}
WorkingDirectory={{ opencode_home }}
Environment=HOME={{ opencode_home }}
EnvironmentFile={{ opencode_env_file }}
ExecStart={{ opencode_npm_prefix.stdout }}/bin/opencode web --hostname {{ opencode_web_hostname }} --port {{ opencode_web_port }}
Restart=always
RestartSec=5s
OOMPolicy=continue
MemoryMax={{ opencode_memory_max }}
MemoryHigh={{ opencode_memory_high }}
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=true
ReadWritePaths={{ opencode_home }}

[Install]
WantedBy=multi-user.target
